# files/systemd/services/bootc-fetch-updates.service
[Unit]
Description=Fetch bootc updates
ConditionPathExists=/run/ostree-booted

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc-fetch-updates
